# deploy/selfhosted/systemd/pseudointelekt-restore-test.service
[Unit]
Description=Verify latest Pseudointelekt backup is readable
After=network-online.target

[Service]
Type=oneshot
User=macie
EnvironmentFile=-/etc/pseudointelekt/backup.env
Environment=RESTIC_REPOSITORY=sftp:restic_pseudointelekt@10.2.11.51:/repo
Environment=RESTIC_PASSWORD_FILE=/etc/pseudointelekt/restic-password
ExecStart=/opt/apps/production/pseudointelekt/deploy/selfhosted/restore-test.sh
PrivateTmp=true
NoNewPrivileges=true
